# files/etc/uci-defaults/99-daed-test-network
#!/bin/sh

# =====================================================================
# 1. 基础网络与管理权精准对齐
# =====================================================================
uci set network.lan.ipaddr='192.168.2.1'
uci set network.lan.netmask='255.255.255.0'
uci set network.lan.gateway='192.168.2.254'
uci -q delete network.lan.dns
uci add_list network.lan.dns='192.168.2.1'
uci set network.lan.delegate='0'
uci commit network

# 配置 SSH (Dropbear) 服务锁定标准 22 端口
uci set dropbear.main.enable='1'
uci set dropbear.main.PasswordAuth='on'
uci set dropbear.main.RootPasswordAuth='on'
uci set dropbear.main.Port='22'
uci set dropbear.main.Interface='lan'
uci commit dropbear

# 强行锁定默认登录密码为 password (双保险机制)
printf "password\npassword\n" | passwd root

# 强制默认主题为 Argon
uci set luci.main.mediaurlbase='/luci-static/argon'
uci commit luci

# =====================================================================
# 2. 🎬 Samba4 局域网 4K 播放多通道吞吐量深度调优
# =====================================================================
SMB_CONF="/etc/config/samba4"
if [ -f "$SMB_CONF" ]; then
    sed -i "/config samba/a \\\tlist server_multi_channel_support 'yes'" $SMB_CONF
    sed -i "/config samba/a \\\tlist rpc_daemon_smbd 'embedded'" $SMB_CONF
    sed -i "/config samba/a \\\tlist aio_read_size '4096'" $SMB_CONF
    sed -i "/config samba/a \\\tlist aio_write_size '4096'" $SMB_CONF
    sed -i "/config samba/a \\\tlist read_raw 'yes'" $SMB_CONF
    sed -i "/config samba/a \\\tlist write_raw 'yes'" $SMB_CONF
    sed -i "/config samba/a \\\tlist use_sendfile 'yes'" $SMB_CONF
fi

# =====================================================================
# 3. ⏰ 自动化运维计划任务注入 (凌晨 4:30 内存释放，5:00 定时重启)
# =====================================================================
mkdir -p /etc/crontabs
cat << 'EOF' >> /etc/crontabs/root
30 4 * * * sync && echo 3 > /proc/sys/vm/drop_caches
0 5 * * * reboot
EOF

# =====================================================================
# 4. 👑 影视聚合后台异步静默构建
# =====================================================================
(
    sleep 60
    echo "=== 开始静默构建 AList 影视聚合 Docker 项目 ==="
    bash -c "$(curl -sSLf https://ailg.ggbond.org/xy_install.sh)"
) &

# 重启服务使配置生效
/etc/init.d/network restart
/etc/init.d/dropbear restart

exit 0
